# povserver/cgi-bin/renderer/loop.cron
PATH = "/usr/local/bin/"
30 * * * * (cd /home/ubuntu/www/cgi-bin/renderer/; ./loop.py)

